# deploy/led-ticker.service
[Unit]
Description=LED Ticker Display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/led-ticker
EnvironmentFile=/opt/led-ticker/.env
ExecStart=/opt/led-ticker/venv/bin/led-ticker --config /opt/led-ticker/config.toml
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# Harden the service. Root is required for GPIO access; the options
# below constrain what a compromised process can do with root.
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
