# deploy/led-ticker-webui.service
# Bare-metal sidecar for the led-ticker web status UI.
# Docker users: the compose.yaml webui service covers this — you do not need
# to install or enable this unit.
#
# Deliberately independent of led-ticker.service: the sidecar must keep
# serving (showing "stale" / "missing") when the display process is down —
# that is exactly when you want the status page. Start order is irrelevant;
# each side degrades gracefully without the other.

[Unit]
Description=LED Ticker Web Status UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# DynamicUser: systemd provisions an ephemeral unprivileged user — no GPIO,
# no useradd step at install. The webui only reads config.toml and
# /run/led-ticker/status.json (written world-readable by the display
# process), and binds an unprivileged port.
DynamicUser=yes
WorkingDirectory=/opt/led-ticker
ExecStart=/opt/led-ticker/venv/bin/led-ticker webui --config /opt/led-ticker/config.toml
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# No .env / EnvironmentFile on purpose: the sidecar needs no API keys, and
# keeping them out of the unprivileged web process preserves the privilege
# separation this sidecar exists for.
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
